# files/etc/hotplug.d/iface/15-mwan3
#!/bin/sh

. /lib/functions.sh
. /lib/functions/network.sh
. /lib/mwan3/common.sh
. /lib/mwan3/mwan3.sh

[ -z "$INTERFACE" ] && exit 0

# Only handle interfaces configured in mwan3
config_load mwan3
enabled=0
config_get enabled "$INTERFACE" enabled 0
[ "$enabled" -eq 1 ] || exit 0

# Load marks from state (hotplug doesn't call mwan3_init)
[ -z "$MMX_MASK" ] && MMX_MASK=$(cat "${MWAN3_STATUS_DIR}/mmx_mask" 2>/dev/null)
[ -z "$MMX_DEFAULT" ] && MMX_DEFAULT=$(cat "${MWAN3_STATUS_DIR}/mmx_default" 2>/dev/null)
[ -z "$MMX_UNREACHABLE" ] && MMX_UNREACHABLE=$(cat "${MWAN3_STATUS_DIR}/mmx_unreachable" 2>/dev/null)
[ -z "$MMX_BLACKHOLE" ] && MMX_BLACKHOLE=$(cat "${MWAN3_STATUS_DIR}/mmx_blackhole" 2>/dev/null)

case "$ACTION" in
	ifup)
		mwan3_lock "hotplug" "$INTERFACE"

		# Remove stale recovery rule (interface back online)
		mwan3_del_recovery_rule "$INTERFACE"

		# Setup routing
		mwan3_create_iface_route "$INTERFACE"
		mwan3_create_iface_rules "$INTERFACE"

		# Setup nft rules for this interface
		mwan3_create_iface_nft "$INTERFACE"

		# Rebuild general nft framework if needed (must be before set population)
		$NFT list chain $MWAN3_NFT_TABLE mwan3_prerouting 2>/dev/null | \
			grep -q "meta mark" || mwan3_set_general_nft

		# Update connected and custom sets (after framework exists)
		mwan3_set_connected_sets
		mwan3_set_custom_sets

		# Rebuild all policies (this interface may be a member)
		mwan3_set_policies_nft

		# Rebuild user rules
		mwan3_set_user_rules

		# Mark interface online and start tracking
		mwan3_set_iface_hotplug_state "$INTERFACE" "online"
		mwan3_track "$INTERFACE"

		mwan3_flush_conntrack "$INTERFACE" "ifup"
		mwan3_unlock "hotplug" "$INTERFACE"
		;;

	ifdown)
		mwan3_lock "hotplug" "$INTERFACE"

		# Kill tracker
		kill $(cat "${MWAN3TRACK_STATUS_DIR}/${INTERFACE}/PID" 2>/dev/null) \
			2>/dev/null

		# Mark offline
		mwan3_set_iface_hotplug_state "$INTERFACE" "offline"

		# Selective conntrack flush for this interface's flows
		mwan3_flush_conntrack "$INTERFACE" "ifdown"

		# WireGuard re-handshake so tunnels move to surviving WAN
		mwan3_wg_rehandshake

		# Remove nft and routing rules
		mwan3_delete_iface_nft "$INTERFACE"
		mwan3_delete_track_iface_nft "$INTERFACE"
		mwan3_delete_iface_rules "$INTERFACE"
		mwan3_delete_iface_route "$INTERFACE"

		# Rebuild policies (this interface removed from pool)
		mwan3_set_policies_nft

		mwan3_unlock "hotplug" "$INTERFACE"
		;;
esac
