# files/etc/hotplug.d/iface/16-mwan3
#!/bin/sh
# mwan3track hotplug handler
# Fired by mwan3track when interface quality changes:
#   connected   → tracking confirms interface is healthy
#   disconnected → tracking confirms interface is down
#   degraded    → interface is lossy but not fully down (future: FEC trigger)

. /lib/functions.sh
. /lib/functions/network.sh
. /lib/mwan3/common.sh
. /lib/mwan3/mwan3.sh

[ -z "$INTERFACE" ] && exit 0
[ -z "$ACTION" ] && exit 0

# Load marks from state
[ -z "$MMX_MASK" ] && MMX_MASK=$(cat "${MWAN3_STATUS_DIR}/mmx_mask" 2>/dev/null)
[ -z "$MMX_DEFAULT" ] && MMX_DEFAULT=$(cat "${MWAN3_STATUS_DIR}/mmx_default" 2>/dev/null)
[ -z "$MMX_UNREACHABLE" ] && MMX_UNREACHABLE=$(cat "${MWAN3_STATUS_DIR}/mmx_unreachable" 2>/dev/null)
[ -z "$MMX_BLACKHOLE" ] && MMX_BLACKHOLE=$(cat "${MWAN3_STATUS_DIR}/mmx_blackhole" 2>/dev/null)

config_load mwan3

case "$ACTION" in
	connected)
		mwan3_lock "hotplug" "$INTERFACE"

		# Interface confirmed healthy by tracker
		# Remove recovery rule — track_hook in nft handles pings now
		mwan3_del_recovery_rule "$INTERFACE"

		mwan3_set_iface_hotplug_state "$INTERFACE" "online"

		# Rebuild policies — interface back in pool
		mwan3_set_policies_nft

		mwan3_flush_conntrack "$INTERFACE" "connected"

		mwan3_unlock "hotplug" "$INTERFACE"
		;;

	disconnected)
		mwan3_lock "hotplug" "$INTERFACE"

		# Interface confirmed down by tracker
		# Add src-IP recovery rule so mwan3track pings can still
		# reach their targets via the interface routing table
		# (track_hook won't fire for offline interfaces since
		#  their chains are removed, so we need the ip rule fallback)
		mwan3_set_recovery_rule "$INTERFACE"

		mwan3_set_iface_hotplug_state "$INTERFACE" "offline"

		# Selective conntrack flush — move existing flows to other WANs
		mwan3_flush_conntrack "$INTERFACE" "disconnected"

		# WireGuard re-handshake
		mwan3_wg_rehandshake

		# Rebuild policies — interface removed from pool
		mwan3_set_policies_nft

		mwan3_unlock "hotplug" "$INTERFACE"
		;;

	degraded)
		# Interface is lossy but not dead
		# Future: trigger UDPSpeeder FEC via FIFO
		# Future: adjust BIRD OSPF cost
		LOG notice "Interface $INTERFACE degraded (loss=${LOSS}% latency=${LATENCY}ms)"

		# Find associated UDPSpeeder tunnel and increase FEC
		for fifo in /tmp/udpspeeder-*${INTERFACE}*.fifo \
		            /tmp/udpspeeder-*$(echo $INTERFACE | sed 's/wan//')*.fifo; do
			[ -p "$fifo" ] && echo "fec 1:3,2:4,8:6,20:10" > "$fifo" 2>/dev/null && \
				LOG notice "Enabled adaptive FEC on $fifo for degraded $INTERFACE"
		done
		;;
esac
