# files/etc/hotplug.d/iface/25-mwan3
#!/bin/sh

. /lib/functions.sh
. /lib/functions/network.sh
. /lib/mwan3/common.sh
. /lib/mwan3/mwan3.sh

[ "$ACTION" = "ifup" ] || [ "$ACTION" = "ifdown" ] || \
[ "$ACTION" = "connected" ] || [ "$ACTION" = "disconnected" ] || \
[ "$ACTION" = "degraded" ] || exit 0

[ -n "$INTERFACE" ] || exit 0

[ -z "$MMX_MASK" ] && MMX_MASK=$(cat "${MWAN3_STATUS_DIR}/mmx_mask" 2>/dev/null)
[ -z "$MMX_DEFAULT" ] && MMX_DEFAULT=$(cat "${MWAN3_STATUS_DIR}/mmx_default" 2>/dev/null)
[ -z "$MMX_UNREACHABLE" ] && MMX_UNREACHABLE=$(cat "${MWAN3_STATUS_DIR}/mmx_unreachable" 2>/dev/null)
[ -z "$MMX_BLACKHOLE" ] && MMX_BLACKHOLE=$(cat "${MWAN3_STATUS_DIR}/mmx_blackhole" 2>/dev/null)

# If state files missing (first boot or /var/run cleared), init now
[ -z "$MMX_MASK" ] && mwan3_init

config_load mwan3

enabled_global=0
config_get_bool enabled_global globals enabled 0
[ "$enabled_global" -gt 0 ] || exit 0

if [ "$ACTION" = "ifup" ] || [ "$ACTION" = "ifdown" ]; then
	iface_enabled=0
	config_get iface_enabled "$INTERFACE" enabled 0
	[ "$iface_enabled" -eq 1 ] || exit 0
fi

# KEY: position 25 runs AFTER 20-firewall (fw4 reload).
# Detect fw4 reload: mwan3_prerouting exists but has no rules.
if [ "$ACTION" = "ifup" ] || [ "$ACTION" = "connected" ]; then
	if ! $NFT list chain $MWAN3_NFT_TABLE mwan3_prerouting 2>/dev/null | \
			grep -q "meta mark"; then
		LOG notice "mwan3 rules missing after fw4 reload, rebuilding"
		mwan3_set_general_nft
		mwan3_set_connected_sets
		mwan3_set_custom_sets
		config_foreach mwan3_rebuild_iface_nft interface
		mwan3_set_policies_nft
		mwan3_set_user_rules
		killall -HUP dnsmasq 2>/dev/null
	fi
fi

case "$ACTION" in
	ifup)
		[ -n "$DEVICE" ] || exit 0
		mwan3_lock "hotplug" "$INTERFACE"

		mwan3_del_recovery_rule "$INTERFACE"
		mwan3_create_iface_route "$INTERFACE"
		mwan3_create_iface_rules "$INTERFACE"
		mwan3_create_iface_nft "$INTERFACE"
		mwan3_update_track_iface_nft "$INTERFACE"
		mwan3_set_connected_sets
		mwan3_set_custom_sets
		mwan3_set_iface_hotplug_state "$INTERFACE" "online"
		mwan3_set_policies_nft
		mwan3_set_user_rules
		mwan3_track "$INTERFACE"
		mwan3_rtmon
		mwan3_flush_conntrack "$INTERFACE" "ifup"

		mwan3_unlock "hotplug" "$INTERFACE"
		;;

	ifdown)
		mwan3_lock "hotplug" "$INTERFACE"

		kill $(cat "${MWAN3TRACK_STATUS_DIR}/${INTERFACE}/PID" 2>/dev/null) \
			2>/dev/null
		mwan3_set_iface_hotplug_state "$INTERFACE" "offline"
		mwan3_flush_conntrack "$INTERFACE" "ifdown"
		mwan3_wg_rehandshake
		mwan3_delete_iface_nft "$INTERFACE"
		mwan3_delete_track_iface_nft "$INTERFACE"
		mwan3_delete_iface_rules "$INTERFACE"
		mwan3_delete_iface_route "$INTERFACE"
		mwan3_set_policies_nft

		mwan3_unlock "hotplug" "$INTERFACE"
		;;

	connected)
		mwan3_lock "hotplug" "$INTERFACE"

		mwan3_del_recovery_rule "$INTERFACE"
		mwan3_set_iface_hotplug_state "$INTERFACE" "online"
		mwan3_create_iface_nft "$INTERFACE"
		mwan3_update_track_iface_nft "$INTERFACE"
		mwan3_set_policies_nft
		mwan3_set_user_rules
		mwan3_rtmon
		mwan3_flush_conntrack "$INTERFACE" "connected"

		mwan3_unlock "hotplug" "$INTERFACE"
		;;

	disconnected)
		mwan3_lock "hotplug" "$INTERFACE"

		mwan3_set_recovery_rule "$INTERFACE"
		mwan3_set_iface_hotplug_state "$INTERFACE" "offline"
		mwan3_flush_conntrack "$INTERFACE" "disconnected"
		mwan3_wg_rehandshake
		mwan3_set_policies_nft

		mwan3_unlock "hotplug" "$INTERFACE"
		;;

	degraded)
		LOG notice "Interface $INTERFACE degraded (loss=${LOSS}% latency=${LATENCY}ms)"
		for fifo in /tmp/udpspeeder-*${INTERFACE}*.fifo; do
			[ -p "$fifo" ] && \
				echo "fec 1:3,2:4,8:6,20:10" > "$fifo" 2>/dev/null && \
				LOG notice "Enabled adaptive FEC on $fifo"
		done
		;;
esac
